# supervisor_1000earths.ini
[program:1000earths]
; command=uwsgi --ini uwsgi_example.ini
directory=/home/markmuetz/Projects/flask-1000earths
command=/home/markmuetz/Projects/flask-1000earths/.env/bin/uwsgi --ini uwsgi_example.ini
environment=PATH=/home/markmuetz/Projects/flask-1000earths/.env/bin
;user=www-data
autostart=False
autorestart=False
